# vvc/service/etc/pydota.cron
MAILTO=""
30 * * * * /home/huangdong/pydota/pydota/service/py_dota_collect.sh
